# supervisord.conf
[supervisord]
nodaemon=false

[program:rsyslog]
command=service rsyslog restart